# s6-overlay/s6-rc.d/opnsense-log-prepare/up
#!/command/execlineb -P

s6-envdir ../opnsense-log/env
importas -i logdir LOGDIR
if { mkdir -p $logdir }
if { chown nobody:nogroup $logdir }
chmod 02755 $logdir