# s6-overlay/s6-rc.d/hosts-add-opnsense/up
#!/command/execlineb -P

s6-envdir ../opnsense/env
multisubstitute {
    importas -i opnsense_lan_address LAN_IPV4_ADDRESS
    importas -i opnsense_hostname HOSTNAME
}
backtick -E opnsense_hostname_short {
    sh -c "echo ${HOSTNAME%%.*}"
}
redirfd -a 1 /etc/hosts
echo "${opnsense_lan_address}\t${opnsense_hostname} ${opnsense_hostname_short}"