# s6-overlay/s6-rc.d/dev-kvm-create/up
#!/command/execlineb -P

# Create kvm and tun device nodes if they don't exist (required --privileged)
define dev_path /dev/kvm
define dev_major 10
define dev_minor 232

ifelse { test -e $dev_path }
{ echo "dev-kvm-create: Device ${dev_path} already exists" }

# if {
#     redirfd -w 1 /sys/fs/cgroup/devices/devices.allow
#     echo "c ${dev_major} ${dev_minor} rwm"
# }
mknod $dev_path c $dev_major $dev_minor
