# s6-overlay/s6-rc.d/dev-kvm-create/up
#!/command/execlineb -P

# Create kvm if it does not exist
define dev_path /dev/kvm
define dev_major 10
define dev_minor 232

ifelse { test -e $dev_path }
{ echo "dev-kvm-create: Device ${dev_path} already exists" }

mknod $dev_path c $dev_major $dev_minor
