# s6-overlay/s6-rc.d/opnsense-log/run
#!/command/execlineb -P

s6-envdir ./env
importas -i logdir LOGDIR
logutil-service $logdir