# s6-overlay/s6-rc.d/opnsense/finish
#!/command/execlineb -P

s6-envdir ./env
importas -i pidfile PIDFILE

if -n { test -e $pidfile }
{ echo "opnsense: PID file not found. Do nothing." }

backtick -E pid { cat $pidfile }

foreground {
    pipeline { echo "system_powerdown\\n" }
    nc localhost 7100
}
wait $pid
