# s6-overlay/s6-rc.d/print-container-usage/up
#!/command/execlineb -P

s6-envdir -f ../opnsense/env
multisubstitute {
    importas -i opnsense_ip_address LAN_IPV4_ADDRESS
    importas -i -N banner_message BANNER_MESSAGE
    importas -i api_key_file API_KEY_FILE
}

echo "\nOPNsense is now operational with the following configuration:\n"

$banner_message

"
Web Interface Access:
- Within the container: http://${opnsense_ip_address}/
- Outside the container: https://<container-name>/

Note: To access from the host, publish port 443 to the host.

API Key:
The API key for the root user is stored at: ${api_key_file}
"
